# wiki_monitor.service
[Unit]
Description=Seeed Wiki Monitor Daemon
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=seeed
WorkingDirectory=/home/seeed/chat
ExecStart=/usr/bin/python3 /home/seeed/chat/monitor_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量
Environment=PYTHONPATH=/home/seeed/chat
Environment=OLLAMA_HOST=127.0.0.1:11434

# 资源限制
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
